# docker/systemd/craftsman.service
[Unit]
Description=Craftsman Agent B Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/admin/hunter-craftsman/craftsman
EnvironmentFile=/home/admin/hunter-craftsman/craftsman/.env
ExecStart=/home/admin/hunter-craftsman/.venv/bin/python -m craftsman.cli serve --host 0.0.0.0 --port 8791
Restart=always
RestartSec=5
User=admin
Group=admin

[Install]
WantedBy=multi-user.target
